# .config/systemd/user/lid-lock.service
[Unit]
Description =Lock on lid close
Before=sleep.target

[Service]
ExecStart=%h/.config/lock.sh
Type=oneshot

[Install]
WantedBy=sleep.target
